# debian/postinst
#!/bin/bash
# postinst script

set -e

PACKAGE="nmexec"
VENV_PATH="/var/lib/nmexec/venv"
VENV_TARBALL="/usr/share/sde_venvs/$PACKAGE.tar.gz"

case "$1" in
    configure)
        # Create nmexec group if it doesn't exist
        if ! getent group nmexec > /dev/null; then
            groupadd --system nmexec
        fi
        # Create nmexec user if it doesn't exist
        if ! getent passwd nmexec > /dev/null; then
            useradd --system --gid nmexec --home /var/lib/nmexec --shell /bin/bash nmexec
        fi

        mkdir -p /var/lib/nmexec
        rm -rf "$VENV_PATH"
        mkdir -p "$VENV_PATH"
        tar -zxf "$VENV_TARBALL" -C "$VENV_PATH"
        rm -f "$VENV_TARBALL"

        # Create shared symlink for yolo9 so other services can access it
        ln -s -f "$VENV_PATH/bin/nmexec" /usr/local/bin/nmexec
        ln -s -f "$VENV_PATH/bin/yolo9_train_dual" /usr/local/bin/yolo9_train_dual
        ln -s -f "$VENV_PATH/bin/yolo9_detect" /usr/local/bin/yolo9_detect
        ln -s -f "$VENV_PATH/bin/yolo9_val_dual" /usr/local/bin/yolo9_val_dual

        chown -R nmexec:nmexec /var/lib/nmexec

        if command -v supervisorctl >/dev/null 2>&1; then
            echo "Using supervisor"
            supervisorctl reread || echo "Could not reread supervisor"
            supervisorctl update || echo "Could not update supervisor"
            supervisorctl restart nmexec || echo "Could not restart nmexec"
            rm -rf /etc/systemd/system/nmexec.service
        else
            echo "Configuring nmexec"
            systemctl daemon-reload
            systemctl enable nmexec.service
            systemctl start nmexec.service
            systemctl restart nmexec.service
            rm -rf /etc/supervisor/conf.d/nmexec.conf
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
