# debian/prerm
#!/bin/bash
# prerm script

set -e

case "$1" in
    upgrade)
        if command -v supervisorctl >/dev/null 2>&1; then
            supervisorctl stop nmexec || true
        else
            systemctl stop nmexec.service || true
        fi
    ;;

    remove|purge)
        if command -v supervisorctl >/dev/null 2>&1; then
            supervisorctl stop nmexec || true
        else
            systemctl stop nmexec.service || true
        fi

        if command -v supervisorctl >/dev/null 2>&1; then
            supervisorctl reread || true
            supervisorctl update || true
        else
            systemctl disable nmexec.service || true
            systemctl daemon-reload || true
        fi
    ;;

    failed-upgrade|abort-install|abort-upgrade|disappear)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
